# templates/ansible-notify.service
[Unit]
After=sshd.service

[Service]
ExecStart=/usr/local/bin/ansible-notify.sh

[Install]
WantedBy=default.target
